# data/xdg-desktop-portal-regolith.service.in
[Unit]
Description=Portal service (Regolith implementation)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.regolith
ExecStart=@LIBEXECDIR@/xdg-desktop-portal-regolith
